# dist/systemd/clefd.service
[Unit]
Description=clefd hotkey daemon
Documentation=https://github.com/nshan651/clef

[Service]
ExecStart=/usr/local/bin/clefd
Restart=on-failure

[Install]
WantedBy=default.target
